# systemd/glchat.service
[Unit]
Description=GL-Chat
After=network.target

[Service]
ExecStart=/home/projects/nodejs/GL-Chat/systemd/start.sh
Restart=always
User=projects
Group=projects
Environment=NODE_ENV=production
WorkingDirectory=/home/projects/nodejs/GL-Chat
LimitNOFILE=4096
IgnoreSIGPIPE=false
KillMode=process

[Install]
WantedBy=multi-user.target
